# systemd/user/offlineimap-oneshot.timer
[Unit]
Description=Offlineimap Query Timer

[Timer]
OnBootSec=3m
OnUnitInactiveSec=15m

[Install]
WantedBy=timers.target
